# scripts/com.thomasbruijns.ictbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thomasbruijns.ictbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/thomasbruijns/projects/ict-futures-bot/scripts/run_monitor.sh</string>
    </array>

    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key><false/>
        <key>Crashed</key><true/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/thomasbruijns/.ict-bot/logs/monitor.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/thomasbruijns/.ict-bot/logs/monitor.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/thomasbruijns/projects/ict-futures-bot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/thomasbruijns/projects/ict-futures-bot/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
